# experiments/rccg/fada-monitor/fada-monitor.service
[Unit]
Description=FADA Process Monitor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/monitor
EnvironmentFile=/etc/fada-monitor/env
ExecStart=/usr/bin/python3 monitor.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/ubuntu/monitor/monitor.log
StandardError=append:/home/ubuntu/monitor/monitor.log

[Install]
WantedBy=multi-user.target
